# examples/minimal/systemd/acp-staging.service
[Unit]
Description=ACP minimal staging stack (Docker Compose)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/root/AI-Control-Plane
ExecStart=/usr/bin/docker compose -f examples/minimal/docker-compose.yml up -d --build
ExecStop=/usr/bin/docker compose -f examples/minimal/docker-compose.yml down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
